# deploy/systemd/simbeamd-demo.service
[Unit]
Description=simbeam demo daemon (headless-browser demo device)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simbeam
Group=simbeam
EnvironmentFile=/etc/simbeam/demo.env
# Chromium insists on a writable HOME even with an explicit user-data-dir.
Environment=HOME=/var/lib/simbeam
# SIMCAST_PAIR_SECRET is read by the daemon itself (keeps the pairing URL
# stable across restarts); the rest of the configuration rides in the args.
ExecStart=/usr/local/bin/simbeamd demo $SIMCAST_DEMO_ARGS
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
StateDirectory=simbeam
WorkingDirectory=/var/lib/simbeam

[Install]
WantedBy=multi-user.target
