# deploy/systemd/simbeam-signal.service
[Unit]
Description=simbeam signaling broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simbeam
Group=simbeam
EnvironmentFile=/etc/simbeam/signal.env
ExecStart=/usr/local/bin/simbeam-signal $SIMCAST_SIGNAL_ARGS
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
StateDirectory=simbeam
WorkingDirectory=/var/lib/simbeam

[Install]
WantedBy=multi-user.target
